# deploy/macos/com.gerenciador.erp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gerenciador.erp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/ruby</string>
        <string>/Users/Shared/GerenciadorERP/app.rb</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/Shared/GerenciadorERP</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>APP_PORT</key>
        <string>4568</string>
        <key>APP_HOST</key>
        <string>0.0.0.0</string>
        <key>RACK_ENV</key>
        <string>production</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/Shared/GerenciadorERP/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/GerenciadorERP/logs/stderr.log</string>
</dict>
</plist>
